# files/etc/systemd/system/btrbk.timer.d/override.conf
[Unit]
Description=btrbk scheduled backup

[Timer]
OnCalendar=*:0/15
AccuracySec=1min
